# setup/etc/init.d/sidekiq_prompt_reports
#!/bin/sh

# chkconfig: 345 85 15
# description: sidekiq_prompt_reports

/etc/init.d/patient_sidekiqctl.rb -e production --dir /var/egotter --user ec2-user --name sidekiq_prompt_reports --state $1
